# files/parent.service
[Unit]
Description=parent

[Service]
ExecStart=/bin/bash -c 'trap "echo exited parent" EXIT; echo started parent; sleep infinity'
